# Source/Data/BrainGenix-EVM.service
[Unit]
Description=BrainGenix EVM Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/BrainGenix-EVM
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target


